# docker/all-in-one/supervisord.conf
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
pidfile=/run/supervisord.pid
logfile=/dev/null

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /data/postgresql -h 127.0.0.1 -p %(ENV_POSTGRES_PORT)s
user=postgres
priority=10
autostart=true
autorestart=unexpected
stopasgroup=true
killasgroup=true
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mysql]
command=/usr/sbin/mysqld --datadir=/data/mysql --socket=/run/mysqld/mysqld.sock --bind-address=127.0.0.1 --port=%(ENV_MYSQL_PORT)s
user=mysql
priority=10
autostart=true
autorestart=unexpected
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mongodb]
command=/usr/bin/mongod --dbpath /data/mongodb --bind_ip 127.0.0.1 --port 27017 --auth --unixSocketPrefix /run/mongodb
user=mongodb
priority=10
autostart=true
autorestart=unexpected
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=/usr/bin/redis-server /run/redis/redis.conf
user=redis
priority=10
autostart=true
autorestart=unexpected
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:api]
command=/usr/local/bin/start-api.sh
directory=/app
user=root
priority=50
autostart=true
autorestart=unexpected
startsecs=10
startretries=3
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
